# scripts/leviathan.service
[Unit]
Description=Leviathan DAPS-Ω
After=network.target

[Service]
User=leviathan
WorkingDirectory=/opt/leviathan
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
EnvironmentFile=/opt/leviathan/.env

[Install]
WantedBy=multi-user.target
